# build/service/webhook/run
#!/bin/bash

php -S localhost:80 -t /usr/src/emby_pinyin
